# fileserve_backend/fileserv.service
[Unit]
Description=FileServ - File Server with Authentication
Documentation=https://github.com/your-org/fileserv
After=network.target

[Service]
Type=simple
# Run as root for PAM authentication and system management
User=root
Group=root
WorkingDirectory=/opt/fileserv
ExecStart=/opt/fileserv/fileserv
Restart=on-failure
RestartSec=5s

# Environment (minimal - most config is in database)
Environment="PORT=8080"
Environment="DATA_DIR=/var/lib/fileserv/data"

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
